# com.smitos.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smitos.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dominium/Library/Application Support/SMIT-OS/start-server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dominium</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/dominium/Library/Logs/SMIT-OS/launchagent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dominium/Library/Logs/SMIT-OS/launchagent-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/sbin</string>
        <key>HOME</key>
        <string>/Users/dominium</string>
    </dict>
</dict>
</plist>
